# daemon/systemd/steam-animation-manager.service
[Unit]
Description=Steam Animation Manager
Documentation=https://github.com/YourUsername/steam-animation-manager
After=graphical-session.target
Wants=graphical-session.target
PartOf=graphical-session.target

[Service]
Type=notify
ExecStart=/usr/bin/steam-animation-daemon --config /etc/steam-animation-manager/config.toml
ExecReload=/bin/kill -HUP $MAINPID
Restart=always
RestartSec=10
TimeoutStopSec=30

# Security settings
User=deck
Group=deck
NoNewPrivileges=yes
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=/tmp /var/tmp /home/deck/.steam /home/deck/.local/share/steam-animation-manager

# Capabilities needed for bind mounts
AmbientCapabilities=CAP_SYS_ADMIN
CapabilityBoundingSet=CAP_SYS_ADMIN

# Resource limits
MemoryMax=256M
CPUQuota=50%

# Environment
Environment=RUST_LOG=info
Environment=RUST_BACKTRACE=1

# Watchdog
WatchdogSec=30
NotifyAccess=main

# Standard streams
StandardOutput=journal
StandardError=journal
SyslogIdentifier=steam-animation-manager

[Install]
WantedBy=graphical-session.target