# bash-daemon/steam-animation-manager.service
[Unit]
Description=Steam Animation Manager (Bash)
Documentation=https://github.com/YourUsername/steam-animation-manager

[Service]
Type=simple
ExecStart=/usr/local/bin/steam-animation-daemon.sh start
ExecReload=/bin/kill -HUP $MAINPID
Restart=always
RestartSec=5
TimeoutStartSec=30
TimeoutStopSec=30

# Environment
Environment=PATH=/usr/local/bin:/usr/bin:/bin

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=steam-animation-manager

[Install]
WantedBy=default.target